# supervisord.conf
[program:apache2]
command=/apache2-start.sh
numprocs=1
autostart=true
autorestart=true
[program:mysqld]
command=/mysqld-start.sh
numprocs=1
autostart=true
autorestart=true
